# files/os-installer-extra/run-mount-squash.mount.in
[Unit]
Description=Mount system squashfs
RequiresMountsFor=/run/mount/installercd/

[Mount]
What=/run/mount/installercd/rlxos.sfs
Where=/run/mount/installercd/squashfs
Type=squashfs
Options=loop

