# files/os-installer-extra/run-mount-squash.mount
[Unit]
Description=Mount rlxos squash

[Mount]
What=/run/mount/installercd/rlxos.sfs
Where=/run/mount/squash
Type=squashfs
Options=defaults
